# Project3/SocketActivation/SocketService-Accepted/SocketService-Accepted.socket
[Unit]
Description=Socket unit for ScoketService

[Socket]
ListenStream=127.0.0.1:4001
Accept=yes
