# Project3/SocketActivation/SocketService-Accepted/SocketService-Accepted@.service
[Unit]
Description=.Net socket activated service

[Service]
ExecStart=/usr/local/bin/dotnet /opt/Project3/SocketService-Accepted.dll
