# Project3/SocketActivation/SocketService-Inet/SocketService-Inet.socket
[Unit]
Description=Socket unit for ScoketService

[Socket]
ListenStream=127.0.0.1:4002
Accept=yes
